# autostart/spirit-tube-tv.desktop
[Desktop Entry]
Type=Application
Name=Spirit Tube TV
Comment=RTL-SDR spirit box — FM, AM, and analog TV bands
Terminal=false
Exec=__INSTALL_DIR__/launch.sh
X-GNOME-Autostart-enabled=true
